# models/Raspberry_Pi/Pi/NXT_SPP_boot
### BEGIN INIT INFO
# Provides: NXT_SPP_boot
# Required-Start: $remote_fs $syslog networking bluetooth
# Required-Stop: $remote_fs $syslog
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: Start serial port server at boot time
# Description: Start serial port server at boot time.
### END INIT INFO

#! /bin/sh

case "$1" in
 start)
   echo "Starting serialPort_serv"
   /home/pi/bluetooth/serialPort_serv &
   sleep 2
   echo "Starting etrobocon_host"
   /home/pi/etrobocon_host_rtt/MW/etrobocon_host &
   ;;

 stop)
   echo "Stopping serialPort_serv"
   pkill serialPort_
   echo "Stopping etrobocon_host" 
   pkill etroboc
   ;;

restart)
   echo "Restart NXT_SPP_boot"
   $0 stop
   sleep 2
   $0 start
   ;;

 *)
   echo "Usage: /etc/init.d/NXT_SPP_boot {start|stop|restart}"
   exit 1
   ;;
esac

exit 0
